# services/ange-gardien/ange-gardien.service
[Unit]
Description=Ange Gardien System Monitor
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/bin/ange-gardien
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target 