# deprecated/spycraft_daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spycraft.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/spycraft</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/spycraft/spycraft.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/spycraft/spycraft.error.log</string>
    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/spycraft</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>1</integer>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
    </dict>
</dict>
</plist> 